# Fundani/Fundani/rpi/udev/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/Fundani/rpi/bin/starter.sh"
